# server/bmw-diag-api.service
[Unit]
Description=BMW Diagnostic AI API Server
After=network.target

[Service]
Type=simple
# Security: Run as dedicated user, not root
User=www-data
Group=www-data
WorkingDirectory=/srv/taller/server
# Load environment from file securely
EnvironmentFile=/srv/taller/server/.env
ExecStart=/usr/bin/node /srv/taller/server/dist/index.js
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/taller/server/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
